# extern/ib1.service
[Unit]
Description=IB1 server
After=network.target

[Service]
StandardOutput=append:/var/log/ib1.log
ExecStart=/usr/local/bin/IB1
WorkingDirectory=/var/lib/ib1
PrivateTmp=true
User=ib
Group=ib

[Install]
WantedBy=multi-user.target
